# agent/pax-agent.service
[Unit]
Description=Aeres PAX push agent (one-shot)
After=network.target

[Service]
Type=oneshot
EnvironmentFile=/etc/aeres/pax-agent.env
ExecStart=/opt/aeres/pax-agent.sh
StandardOutput=journal
StandardError=journal
User=root

[Install]
WantedBy=multi-user.target
